# systemd/eg4-refresh-report.timer
[Unit]
Description=Refresh Solar Digital Twin EG4 reports every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true
Unit=eg4-refresh-report.service

[Install]
WantedBy=timers.target
